# db-pull-on-boot/templates/db-pull.service.j2
[Unit]
Description=DB Pull Service

[Service]
type=oneshot
ExecStart=/home/ubuntu/db-push.sh
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
