# fix-hosts-file-on-boot/templates/fix-hosts.service.j2
[Unit]
Description=Fix Hosts service

[Service]
Type=oneshot
ExecStart={{ dest_location }}/fix-hosts.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
