# fix-hosts-file-on-boot/templates/fix-hosts.service.j2
[Unit]
Description=Fix Hosts service

[Service]
Type=oneshot
ExecStart={{ chia_root }}/fix-hosts.sh
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
